# var/lib/dpkg/info/cups-daemon.list
/.
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/cups.socket
/lib/systemd/system/cups.service
/lib/systemd/system/cups.path
/usr
/usr/lib
/usr/lib/cups
/usr/lib/cups/notifier
/usr/lib/cups/notifier/mailto
/usr/lib/cups/notifier/dbus
/usr/lib/cups/notifier/rss
/usr/lib/cups/backend
/usr/lib/cups/backend/ipp
/usr/lib/cups/backend/ipp14
/usr/share
/usr/share/doc
/usr/share/doc/cups-daemon
/usr/share/doc/cups-daemon/examples
/usr/share/doc/cups-daemon/examples/cups-socket.public.conf
/usr/share/doc/cups-daemon/examples/cups-socket.localhost.conf
/usr/share/doc/cups-daemon/NEWS.Debian.gz
/usr/share/doc/cups-daemon/copyright
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_cups.py
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/cups-daemon
/usr/share/man
/usr/share/man/man7
/usr/share/man/man7/backend.7.gz
/usr/share/man/man8
/usr/share/man/man8/cupsd.8.gz
/usr/share/man/man5
/usr/share/man/man5/cupsd.conf.5.gz
/usr/share/man/man5/mailto.conf.5.gz
/usr/share/man/man5/printers.conf.5.gz
/usr/share/man/man5/mime.types.5.gz
/usr/share/man/man5/cups-snmp.conf.5.gz
/usr/share/man/man5/cups-files.conf.5.gz
/usr/share/man/man5/classes.conf.5.gz
/usr/share/cups
/usr/share/cups/cupsd.conf.default
/usr/share/cups/mime
/usr/share/cups/mime/mime.types
/usr/sbin
/usr/sbin/cupsd
/etc
/etc/init.d
/etc/init.d/cups
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.cupsd
/etc/pam.d
/etc/pam.d/cups
/etc/init
/etc/init/cups.conf
/etc/logrotate.d
/etc/logrotate.d/cups-daemon
/etc/ufw
/etc/ufw/applications.d
/etc/ufw/applications.d/cups
/etc/cups
/etc/cups/ssl
/etc/cups/cups-files.conf
/var
/var/spool
/var/spool/cups
/var/spool/cups/tmp
/var/log
/var/log/cups
/var/cache
/var/cache/cups
/var/cache/cups/rss
/usr/lib/cups/backend/https
/usr/lib/cups/backend/ipps
/usr/lib/cups/backend/http
/usr/share/doc/cups-daemon/README.Debian.gz
/usr/share/doc/cups-daemon/changelog.Debian.gz
