# var/lib/dpkg/info/cups-browsed.list
/.
/usr
/usr/sbin
/usr/sbin/cups-browsed
/usr/share
/usr/share/man
/usr/share/man/man8
/usr/share/man/man8/cups-browsed.8.gz
/usr/share/man/man5
/usr/share/man/man5/cups-browsed.conf.5.gz
/usr/share/doc
/usr/share/doc/cups-browsed
/usr/share/doc/cups-browsed/changelog.Debian.gz
/usr/share/doc/cups-browsed/README.gz
/usr/share/doc/cups-browsed/copyright
/usr/share/doc/cups-browsed/AUTHORS
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/cups-browsed.service
/etc
/etc/init
/etc/init/cups-browsed.conf
/etc/cups
/etc/cups/cups-browsed.conf
/etc/apparmor.d
/etc/apparmor.d/usr.sbin.cups-browsed
/etc/init.d
/etc/init.d/cups-browsed
