# var/lib/dpkg/info/kerneloops-daemon.conffiles
/etc/dbus-1/system.d/kerneloops.dbus
/etc/init.d/kerneloops
/etc/kerneloops.conf
/etc/default/kerneloops
