# var/lib/dpkg/info/plymouth.list
/.
/var
/var/spool
/var/spool/plymouth
/var/lib
/var/lib/plymouth
/etc
/etc/init
/etc/init/plymouth-upstart-bridge.conf
/etc/init/plymouth-stop.conf
/etc/init/plymouth-log.conf
/etc/init/plymouth-splash.conf
/etc/init/plymouth-ready.conf
/etc/init/plymouth.conf
/etc/init/plymouth-shutdown.conf
/lib
/lib/systemd
/lib/systemd/system
/lib/systemd/system/kexec.target.wants
/lib/systemd/system/plymouth-poweroff.service
/lib/systemd/system/poweroff.target.wants
/lib/systemd/system/plymouth-switch-root.service
/lib/systemd/system/plymouth-read-write.service
/lib/systemd/system/plymouth-start.service
/lib/systemd/system/plymouth-kexec.service
/lib/systemd/system/plymouth-halt.service
/lib/systemd/system/halt.target.wants
/lib/systemd/system/reboot.target.wants
/lib/systemd/system/sysinit.target.wants
/lib/systemd/system/plymouth-quit-wait.service
/lib/systemd/system/systemd-ask-password-plymouth.path
/lib/systemd/system/plymouth-reboot.service
/lib/systemd/system/initrd-switch-root.target.wants
/lib/systemd/system/multi-user.target.wants
/lib/systemd/system/systemd-ask-password-plymouth.service
/lib/systemd/system/plymouth-quit.service
/lib/plymouth
/lib/plymouth/ubuntu_logo.png
/lib/plymouth/themes
/lib/plymouth/themes/details
/lib/plymouth/themes/details/details.plymouth
/lib/arm-linux-gnueabihf
/lib/arm-linux-gnueabihf/plymouth
/lib/arm-linux-gnueabihf/plymouth/script.so
/lib/arm-linux-gnueabihf/plymouth/details.so
/lib/arm-linux-gnueabihf/plymouth/renderers
/lib/arm-linux-gnueabihf/plymouth/renderers/frame-buffer.so
/lib/arm-linux-gnueabihf/plymouth/renderers/drm.so
/bin
/bin/plymouth
/bin/plymouth-upstart-bridge
/usr
/usr/share
/usr/share/initramfs-tools
/usr/share/initramfs-tools/hooks
/usr/share/initramfs-tools/hooks/plymouth
/usr/share/initramfs-tools/scripts
/usr/share/initramfs-tools/scripts/panic
/usr/share/initramfs-tools/scripts/panic/plymouth
/usr/share/initramfs-tools/scripts/init-top
/usr/share/initramfs-tools/scripts/init-top/plymouth
/usr/share/initramfs-tools/scripts/init-bottom
/usr/share/initramfs-tools/scripts/init-bottom/plymouth
/usr/share/initramfs-tools/conf-hooks.d
/usr/share/apport
/usr/share/apport/package-hooks
/usr/share/apport/package-hooks/source_plymouth.py
/usr/share/doc
/usr/share/doc/plymouth
/usr/share/doc/plymouth/copyright
/usr/share/doc/plymouth/README.Debian
/sbin
/sbin/plymouthd
/lib/systemd/system/kexec.target.wants/plymouth-kexec.service
/lib/systemd/system/poweroff.target.wants/plymouth-poweroff.service
/lib/systemd/system/halt.target.wants/plymouth-halt.service
/lib/systemd/system/reboot.target.wants/plymouth-reboot.service
/lib/systemd/system/sysinit.target.wants/plymouth-read-write.service
/lib/systemd/system/sysinit.target.wants/plymouth-start.service
/lib/systemd/system/initrd-switch-root.target.wants/plymouth-switch-root.service
/lib/systemd/system/initrd-switch-root.target.wants/plymouth-start.service
/lib/systemd/system/multi-user.target.wants/plymouth-quit-wait.service
/lib/systemd/system/multi-user.target.wants/plymouth-quit.service
/usr/share/doc/plymouth/changelog.Debian.gz
